# etc/units/example.timer
[Timer]
OnBootSec=10s
OnUnitActiveSec=1h
Unit=example.service
